# rootfs/etc/services.d/ecoline-telegram-bot/run
#!/usr/bin/with-contenv sh
cd /opt/ecoline-telegram-bot
if [ -z "$ECOLINE_DEBUG" ]
then

    python /opt/ecoline-telegram-bot/bot.py --debug
else
    python /opt/ecoline-telegram-bot/bot.py
fi
